# scripts/cron.d/claude-sandbox-cleanup
# Cleanup orphaned worktrees every hour
0 * * * * claude-sandbox /srv/claude-sandbox/app/scripts/cleanup-worktrees.sh
